# etc/systemd/system/doas_startup.service
[Unit]
Description= doas startup
After=hibernate.target

[Service]
ExecStart=/etc/doas_startup

[Install]
WantedBy=hibernate.target
